# docker/config/supervisord/kr8851queue.conf
[program:kr8851queue-csvsite]
process_name=%(program_name)s_%(process_num)02d
user=root
command=php /home/website/kr8851queue_web/artisan queue:work --queue=csvsite --tries=5 --memory=512 --sleep=5
autostart=true
autorestart=true
stderr_logfile=/home/logs/supervisord/kr8851queue/csvsite/error.log
stdout_logfile=/home/logs/supervisord/kr8851queue/csvsite/nohup.log


[program:kr8851queue-olddata]
process_name=%(program_name)s_%(process_num)02d
user=root
command=php /home/website/kr8851queue_web/artisan queue:work --queue=olddata --tries=5 --memory=512 --sleep=5
autostart=true
autorestart=true
stderr_logfile=/home/logs/supervisord/kr8851queue/olddata/error.log
stdout_logfile=/home/logs/supervisord/kr8851queue/olddata/nohup.log


[program:kr8851queue-official]
process_name=%(program_name)s_%(process_num)02d
user=root
command=php /home/website/kr8851queue_web/artisan queue:work --queue=official --tries=5 --memory=512 --sleep=5
autostart=true
autorestart=true
stderr_logfile=/home/logs/supervisord/kr8851queue/official/error.log
stdout_logfile=/home/logs/supervisord/kr8851queue/official/nohup.log

